# lib/template/rackspace-onmetal-io/lsi-settings.service
[Unit]
Description=Configure performance settings for LSI cards

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/root/bin/lsi_settings.sh