# lib/template/format-ephemeral.service
[Unit]
Description=Formats the ephemeral drive

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/root/bin/format-ephemeral.sh